# .deploy/prod/supervisor/new-arthall.octane-worker.ini
[program:new-arthall-octane]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/arthall/arthall.newback/artisan octane:start --max-requests=250 --workers=1 --port=8091
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/arthall/arthall.newback/storage/logs/laravel-octane-worker.log​
stopwaitsecs=3600
stopasgroup=true
stopsignal=KILL
killasgroup=true
